# config/com.dotfiles.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- ジョブのラベル（一意の識別子） -->
    <key>Label</key>
    <string>com.dotfiles.sync</string>

    <!-- 実行するプログラム -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mba/01-dev/dotfiles/scripts/sync_dotfiles.sh</string>
        <string>--with-brew</string>
    </array>

    <!-- 定期実行の間隔（秒単位：3600 = 1時間） -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- ログイン時に実行 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 標準出力のログファイル -->
    <key>StandardOutPath</key>
    <string>/Users/mba/Library/Logs/dotfiles-sync.log</string>

    <!-- 標準エラー出力のログファイル -->
    <key>StandardErrorPath</key>
    <string>/Users/mba/Library/Logs/dotfiles-sync-error.log</string>

    <!-- 環境変数 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
